# deploy/vps/systemd/itmcopy-cron.timer
[Unit]
Description=itmCopy — імпорт новин кожні 6 годин

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
